# config/systemd/user/gotopath.service
[Unit]
Description=Go to path

[Service]
ExecStart=#GOTOPATH# -serve

[Install]
WantedBy=default.target